# couchbase/3.0.1/config/supervisor/couchbase.conf
[supervisord]
nodaemon=true

[program:couchbase]
command=/usr/bin/start-couchbase.sh
autostart=true
