# p_learning/1.0/config/supervisor/apache.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/bin/start-apache.sh
autostart=true
