# p_learning/1.0/config/supervisor/mysql.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/start-mysql.sh
autostart=true
